# system/jarvisos-generator.service
[Unit]
Description=JarvisOS Generator - Auto Code Generation
After=jarvisos-analyzer.service
Requires=jarvisos-analyzer.service
Documentation=https://github.com/yourusername/jarvisos

[Service]
Type=oneshot
User=jarvis
Group=jarvis
WorkingDirectory=/opt/jarvisos
Environment="PATH=/opt/jarvisos/venv/bin:/usr/local/bin:/usr/bin"
Environment="ANTHROPIC_API_KEY=your_key_here"

# Run generator
ExecStart=/opt/jarvisos/venv/bin/python /opt/jarvisos/jarvis.py generate

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=jarvisos-generator
